# packaging/deb/prerm
#!/bin/sh
set -eu

if [ "${1:-}" = remove ]; then
    /usr/bin/systemctl disable --now g7tg-agent.service >/dev/null 2>&1 || true
fi

